# ssh/.local/zsh/profile.d/ssh_aliases.zsh
alias ssh-noctrl='\ssh -o ControlMaster=no -o ControlPath=none'
alias ssh-yolo='\ssh -o ControlMaster=no -o ControlPath=none'
