# ssh/.local/zsh/profile.d/ssh_aliases.zsh
alias ssh-gen-config='\assemble-config -v -f ~/.ssh/config -V sshconfig'
alias ssh-noctrl='\ssh -o ControlMaster=no -o ControlPath=none'
alias ssh-yolo='\ssh -o ControlMaster=no -o ControlPath=none'
